# amazon-crawler.service
[Unit]
Description=amazon-crawler — Amazon product data API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/amazon-crawler
Environment=HOME=/root
Environment=DATA_DIR=/opt/amazon-crawler/data
Environment=CAMEL_PROFILE_DIR=/tmp/amazon-crawler-profile
# Different CDP port from sell's ebay-scraper (which uses 9250) so the two
# services' Chromium subprocesses don't collide on the same host.
Environment=CAMEL_CDP_PORT=9260
Environment=CAMEL_XVFB_DISPLAY=:99
ExecStartPre=/bin/sh -c 'rm -f /tmp/amazon-crawler-profile/Singleton*'
ExecStart=/opt/amazon-crawler/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8011
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
