# applesauce-crawlers.service
[Unit]
Description=applesauce-crawlers — multi-source product data API (CCC, eBay, ...)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/applesauce-crawlers
Environment=HOME=/root
Environment=DATA_DIR=/opt/applesauce-crawlers/data
Environment=CAMEL_PROFILE_DIR=/tmp/applesauce-crawlers-camel-profile
# CDP port + Xvfb display picked to NOT collide with the other Chromium
# users on this host:
#   ebay-scraper.service  (sell's CCC scraper) — CDP 9250 / display :98
#   tiktok-scraper.service                     — display :99
Environment=CAMEL_CDP_PORT=9260
Environment=CAMEL_XVFB_DISPLAY=:96
ExecStartPre=/bin/sh -c 'rm -f /tmp/applesauce-crawlers-camel-profile/Singleton* 2>/dev/null || true'
ExecStart=/opt/applesauce-crawlers/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8014
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
